# ops/systemd/iic-morning.timer
[Unit]
Description=IIC-FORGE morning digest timer
Requires=iic-morning.service

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true
Unit=iic-morning.service

[Install]
WantedBy=timers.target
